# app_config.conf
[program:myserver]

command=gunicorn app:BidMeApp -w 4 -b 0.0.0.0:80

[supervisord]
logfile=./bidmeapp_supervisord.log
loglevel=debug
user=root
